# deploy/examples/vm/systemd/ati-rotate-keyring.service
[Unit]
Description=Rotate the ATI keyring from 1Password and SIGHUP the live proxy
Documentation=https://github.com/Parcha-ai/ati/blob/main/deploy/examples/vm/README.md
# No manual ConditionPathExists=! lockfile — that would silently SKIP
# the unit forever if a crash mid-rotation orphaned the lock. Instead
# we use RuntimeDirectory= below: systemd creates it at unit start and
# auto-cleans it on unit stop (including crash). If two rotations
# overlap, the second one fails its mkdir cleanly with EEXIST and
# systemd logs it — far better than a silent skip.

[Service]
Type=oneshot
User=ati
Group=ati

# systemd-managed lock directory under /run/ati-rotate-keyring/. Created
# at unit start, cleaned at unit stop (even after a crash). External
# observers can read this dir's existence as the "rotation in progress"
# signal.
RuntimeDirectory=ati-rotate-keyring
RuntimeDirectoryMode=0700

# 1Password service-account token. Provisioned as root:ati 0440 outside
# this unit so the service can read it via systemd's LoadCredential=
# (which copies the file into a private $CREDENTIALS_DIRECTORY for this
# service alone — the unit's User= can read it without needing access to
# /etc/op-service-account-token directly).
#
# IMPORTANT: %d expands to $CREDENTIALS_DIRECTORY (a PATH), NOT the token
# value. We pass that path to `ati edge rotate-keyring` via --op-token-file,
# which reads the contents and exports them as OP_SERVICE_ACCOUNT_TOKEN to
# the `op` subprocess.
LoadCredential=op-token:/etc/op-service-account-token

# Vault / item to pull. Override per environment.
Environment=ATI_OP_VAULT=Production Secrets
Environment=ATI_OP_ITEM=ATI Edge VM Keyring

ExecStart=/usr/local/bin/ati edge rotate-keyring \
    --vault "${ATI_OP_VAULT}" \
    --item "${ATI_OP_ITEM}" \
    --ati-dir /var/lib/ati \
    --op-token-file %d/op-token \
    --service ati

# Hardening — same posture as ati.service.
ProtectSystem=strict
ReadWritePaths=/var/lib/ati
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
